# files/system/usr/lib/systemd/system/disable-power-wakeups.timer
[Unit]
Description=Disable power wakeups Trigger

[Timer]
RandomizedDelaySec=1m
OnStartupSec=1m

[Install]
WantedBy=timers.target
